# argus-insight-agent/packaging/systemd/argus-agent.service
[Unit]
Description=Argus Insight Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/argus-agent
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argus-agent

# Directories
RuntimeDirectory=argus-insight-agent
StateDirectory=argus-insight-agent
LogsDirectory=argus-insight-agent
ConfigurationDirectory=argus-insight-agent

[Install]
WantedBy=multi-user.target
