# argus-insight-agent/packaging/debian/postinst
#!/bin/sh
set -e

#DEBHELPER#

if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable argus-agent.service
fi
