# argus-insight-agent/packaging/debian/postinst
#!/bin/sh
set -e

#DEBHELPER#

if [ "$1" = "configure" ]; then
    # Create certificate directories
    mkdir -p /opt/argus-insight-agent/certificates/ca

    systemctl daemon-reload
    systemctl enable argus-insight-agent.service
fi
